# files/system/etc/skel/.config/autostart/first-run-sound.desktop
[Desktop Entry]
Type=Application
Name=First Run Sound
Exec=/usr/bin/first-run-sound.sh
X-KDE-autostart-after=panel
OnlyShowIn=KDE;
NoDisplay=true
