# files/system/etc/skel/.config/systemd/user/first-run-sound.service
[Unit]
Description=First run sound

[Service]
Type=oneshot
ExecStart=/usr/bin/first-run-sound.sh

[Install]
WantedBy=default.target
